# nginx/certbot/certbot.service
[Unit]
Description=Renew Certbot certificates
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew
