# nginx/certbot/certbot.timer
[Unit]
Description=Renew Certbot certificate

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=86400

[Install]
WantedBy=multi-user.target
